# templates/systemd/system/daemontools_service.service
[Unit]
Description=Daemontools service supervision

[Service]
ExecStart=/usr/bin/svscanboot {{ daemontools_service_dir }}
Restart=always

[Install]
WantedBy=multi-user.target
