# templates/systemd/system/svscan.service
[Unit]
Description=Daemontools service supervision

[Service]
ExecStart={{ daemontools_dir }}/svscanboot {{ daemontools_svscan_dir }}
Restart=always

[Install]
WantedBy=multi-user.target
